# Library/LaunchAgents/com.orphee.obsidian-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>Label</key>
        <string>com.orphee.obsidian-update</string>
        <key>ProgramArguments</key>
        <array>
                <string>/Users/roger/.config/obsidian/download-icloud.sh</string>
        </array>
        <key>StandardErrorPath</key>
                <string>/tmp/com.orphee.obsidian-update.stderr</string>
        <key>StandardOutPath</key>
                <string>/tmp/com.orphee.obsidian-update.stdout</string>
        <key>StartCalendarInterval</key>
        <array>
                <dict>
                        <key>Weekday</key>
                        <integer>2</integer>
                        <key>Hour</key>
                        <integer>06</integer>
                        <key>Minute</key>
                        <integer>00</integer>
                </dict>
                <dict>
                        <key>Weekday</key>
                        <integer>2</integer>
                        <key>Hour</key>
                        <integer>08</integer>
                        <key>Minute</key>
                        <integer>00</integer>
                </dict>
                <dict>
                        <key>Weekday</key>
                        <integer>2</integer>
                        <key>Hour</key>
                        <integer>09</integer>
                        <key>Minute</key>
                        <integer>30</integer>
                </dict>
                <dict>
                        <key>Weekday</key>
                        <integer>4</integer>
                        <key>Hour</key>
                        <integer>06</integer>
                        <key>Minute</key>
                        <integer>00</integer>
                </dict>
                <dict>
                        <key>Weekday</key>
                        <integer>4</integer>
                        <key>Hour</key>
                        <integer>08</integer>
                        <key>Minute</key>
                        <integer>00</integer>
                </dict>
                <dict>
                        <key>Weekday</key>
                        <integer>4</integer>
                        <key>Hour</key>
                        <integer>09</integer>
                        <key>Minute</key>
                        <integer>30</integer>
                </dict>
                <dict>
                        <key>Weekday</key>
                        <integer>1</integer>
                        <key>Hour</key>
                        <integer>10</integer>
                        <key>Minute</key>
                        <integer>00</integer>
                </dict>
                <dict>
                        <key>Weekday</key>
                        <integer>1</integer>
                        <key>Hour</key>
                        <integer>12</integer>
                        <key>Minute</key>
                        <integer>00</integer>
                </dict>
                <dict>
                        <key>Weekday</key>
                        <integer>5</integer>
                        <key>Hour</key>
                        <integer>08</integer>
                        <key>Minute</key>
                        <integer>00</integer>
                </dict>
                <dict>
                        <key>Weekday</key>
                        <integer>5</integer>
                        <key>Hour</key>
                        <integer>09</integer>
                        <key>Minute</key>
                        <integer>00</integer>
                </dict>
                <dict>
                        <key>Weekday</key>
                        <integer>7</integer>
                        <key>Hour</key>
                        <integer>06</integer>
                        <key>Minute</key>
                        <integer>00</integer>
                </dict>
        </array>
</dict>
</plist>
